# examples/ruby/server_infra/docker/base/build/service/greeter_server/run
#!/bin/bash
cd /home/app/greeter_server
export GRPC_PORT=50061
exec bundle exec ruby greeter_server.rb
