# examples/ruby/server_infra/docker/base/build/service/greeter_server/run
#!/bin/bash
cd /home/app/greeter_server
PORT=50061 bundle exec ruby greeter_server.rb
